# ops/com.you.voicebot-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.you.voicebot-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/rsync</string>
        <string>-rt</string>
        <string>--no-perms</string>
        <string>--no-owner</string>
        <string>--no-group</string>
        <string>--modify-window=2</string>
        <string>/Users/corbin/PycharmProjects/agi_env_v1/ucs2/data/</string>
        <string>/Users/corbin/Library/CloudStorage/GoogleDrive-c@c42.io/My Drive/ucs-backup/</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/corbin/Library/Logs/voicebot/backup-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/corbin/Library/Logs/voicebot/backup-stderr.log</string>
</dict>
</plist>
